# InstallationPackage/kaltura-installer/package/app/app/auto_install/package/app/crontab/kaltura_rotate
30 4 * * * root /usr/sbin/logrotate -f -vv @BASE_DIR@/logrotate/kaltura_log_rotate >> @LOG_DIR@/kaltura_log_rotate.log >> @LOG_DIR@/cron.log 2>&1
